# salesforce-mockup/salesforce-mockup.service
# Systemd unit para Meubles RD Salesforce mockup (uvicorn).
# Instalación (como root o con sudo):
#   sudo cp /opt/gac-copilot/salesforce-mockup/salesforce-mockup.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable salesforce-mockup
#   sudo systemctl start salesforce-mockup
# Ver logs: sudo journalctl -u salesforce-mockup -f

[Unit]
Description=Meubles RD Salesforce mockup (uvicorn)
After=network.target

[Service]
Type=simple
User=rarme
Group=rarme
WorkingDirectory=/opt/gac-copilot/salesforce-mockup
ExecStart=/opt/gac-copilot/salesforce-mockup/start-prod.sh
Restart=always
RestartSec=3
# Puerto configurado en salesforce-mockup/.env (UVICORN_PORT)
# Environment="ALLOWED_HOSTS=localhost,127.0.0.1,tudominio.com"
# Environment="DJANGO_SECRET_KEY=tu-secret-key"

[Install]
WantedBy=multi-user.target
